# uccp420wlan/files/etc/init.d/uccp420wlan
#!/bin/sh /etc/rc.common

START=10
STOP=10

start() {
	modprobe mac80211
	insmod uccp420wlan
}

stop() {
	rmmod uccp420wlan
}
